# Htroy/supervisor/conf/lkl_bbr.ini
[program:lkl_bbr]
command=bash /home/tcp_nanqinlang/load.sh
autorestart=true
autostart=true
user=root
redirect_stderr=true
