# Htroy/supervisor/conf/ocserv.ini
[program:ocserv]
command=ocserv -f -c /etc/ocserv/ocserv.conf
autorestart=true
autostart=true
user=root
redirect_stderr=true
