# Htroy/supervisor/conf/shadowsocks.ini
[program:shadowsocks]
command=python /usr/local/shadowsocks/server.py
autorestart=true
autostart=true
user=shadowsocks
redirect_stderr=true
