# Htroy/supervisor/conf/sslocal.ini
[program:shadowsocks]
command=sslocal -c /etc/shadowsocks/serverd.json
autostart=true
autorestart=true
user=shadowsocks
redirect_stderr=true
