# Htroy/supervisor/conf/udp2rawc.ini
[program:udp2raw2]
command=/usr/bin/udp2raw -c -l 0.0.0.0: -r :9100 -k "" --raw-mode faketcp -a
autorestart=true
autostart=true
user=root
redirect_stderr=true
