# Htroy/supervisor/conf/udp2raw.ini
[program:udp2raw4]
command=/usr/bin/udp2raw -s -l 0.0.0.0:9100 -r 127.0.0.1:9100 -k "" --raw-mode faketcp -a
autorestart=true
autostart=true
user=root
redirect_stderr=true
